# src/installation/sqslistener.service
[Unit]
Description=sqs listener

[Service]
Type=simple
ExecStart=/usr/bin/sqslistener

[Install]
WantedBy=multi-user.target